# sccomans
estos son unos comandos para crear script en debian
para poder escribir un script se debe crear un archivo name.sh
tambien se pueden identificar con el encabezado:
#!/bin/bash ; #!/bin/ksh
para python es:
	#!/usr/bin/env python
el #! recibe le nombre de shebang que inician los programas ejecutables al comienzo.
al guardarlo en el directorio que deseas (la opcion mas comun es guardarlo en /home/nombre_usuario). Existen varios modos de ejecutarlo al inicio y uno de ellos es mover el archivo a la carpeta init.d con el comando:
	sudo mv /home/tu_usuario/nombre_script.sh /etc/init.d/
luego es necesario darle permiso de ejecucion
	sudo chmod +x /etc/init.d/nombre_script.sh
luego actualizaremos el rc.d
	sudo update-rc.d nombre_script.sh defaults
si queremos desintalar solo ejecutamos:
	rm /etc/init.d/nombre_script.sh
Ahora si en el caso de que el comando update-rc.d te muestre el siguiente mensaje:
	insserv: warning: script 'name_script' missing LSB tags and overrides
se tendra que agregar en el scritp:
	#!/bin/bash
	### BEGIN INIT INFO
	# Provides:          blabla
	# Required-Start:    $syslog
	# Required-Stop:     $syslog
	# Default-Start:     2 3 4 5
	# Default-Stop:      0 1 6
	# Short-Description: blabla
	# Description:
	#
	### END INIT INFO
	
	##pag lcd rpi
	
http://blog.bricogeek.com/noticias/electronica/lcd-grafico-con-raspberry-pi-y-python/
